# root/etc/s6-overlay/s6-rc.d/svc-mod-universal-postgres/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid postgres \
	pg_ctl -D /config/postgres/ \
	-o '-c listen_addresses=localhost -p 5432' \
	-w start
