# root/etc/s6-overlay/s6-rc.d/svc-mod-universal-redis/run
#!/usr/bin/with-contenv bash

echo "Starting Redis."
exec s6-setuidgid abc \
	redis-server \
	/config/redis/redis.conf
