# root/etc/s6-overlay/s6-rc.d/init-mod-universal-redis-add-package/run
#!/usr/bin/with-contenv bash

echo "\
    redis-server" >>/mod-pip-packages-to-install.list
